# config/web/nginx/etc/service/agent/run
#!/bin/bash
# Vhost satan agent
exec setuidgid satan /usr/bin/perl /home/satan/prod/agent.pl vhost
