# config/web/mysql/etc/service/agent/run
#!/bin/bash
# MySQL satan agent
exec setuidgid satan /usr/bin/perl /home/satan/prod/agent.pl mysql
