# config/web/nginx/etc/service/nginx/run
#!/bin/bash
# Nginx web server
exec /usr/local/nginx/sbin/nginx -c /etc/nginx/nginx.conf
